# scripts/rootfs-services/tenbox-agent-backup.service
[Unit]
Description=Protect TenBox Agent data

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tenbox-agent-backup snapshot
