# scripts/rootfs-services/virtiofs-automount.service
[Unit]
Description=Auto-mount virtiofs shared folders
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/virtiofs-automount
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
